# systemd/boiler-alarm.service
[Unit]
Description=Boiler Alarm Monitoring Service
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/boiler-alarm
ExecStart=/usr/bin/python3 /home/pi/boiler-alarm/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (alternatively, use EnvironmentFile)
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target